# hedera-node/infrastructure/docker/containers/production-next/consensus-node/etc/s6-overlay/s6-rc.d/hedera/run
#! /command/with-contenv bash

s6-svc -O /run/service/hedera
/opt/hgcapp/services-hedera/HapiApp2.0/entrypoint.sh
